# deploy/marcel-redeploy.service.tmpl
# Marcel — redeploy service (triggered by marcel-redeploy.path)
# Runs redeploy.sh on the host when the agent requests a restart.
# This unit is NOT started manually — the path unit triggers it.
[Unit]
Description=Marcel redeploy (rebuild + restart with rollback)

[Service]
Type=oneshot
WorkingDirectory=@@REPO_DIR@@
ExecStart=@@REPO_DIR@@/redeploy.sh
TimeoutStartSec=300
