# deploy/marcel.service.tmpl
# Marcel — main Docker service (user-level systemd unit)
# Installed by: make setup
# Managed by: systemctl --user {start|stop|restart|status} marcel
[Unit]
Description=Marcel personal agent (Docker)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=@@REPO_DIR@@
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=180

[Install]
WantedBy=default.target
