# deploy/marcel-redeploy.service.tmpl
[Unit]
Description=Rebuild and restart Marcel (prod) after self-modification
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=@@REPO_ROOT@@
ExecStart=@@REPO_ROOT@@/scripts/redeploy.sh --env prod --force
StandardOutput=journal
StandardError=journal
