# deploy/marcel-redeploy.path.tmpl
[Unit]
Description=Watch for Marcel restart request flag

[Path]
PathExists=@@DATA_DIR@@/watchdog/restart_requested
Unit=marcel-redeploy.service

[Install]
WantedBy=default.target
